# config/systemd/user/battery-monitor.service
[Unit]
Description=Omarchy Battery Monitor Check
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/battery-monitor
KillMode=process
Environment=DISPLAY=:0
